# iso_root/etc/systemd/system/serial-getty@ttyS0.service.d/autologin.conf
# /etc/systemd/system/serial-getty@ttyS0.service.d/autologin.conf
# Drop-in override for the systemd serial getty service on ttyS0.
# Configures automatic root login on the serial console so that the
# QEMU smoke-test (run_iso.sh --wait-boot) can detect the login prompt
# in the serial log and confirm a successful boot.

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --keep-baud 115200,57600,38400,9600 %I $TERM
