# iso_root/etc/systemd/system/systemd-networkd-wait-online.service.d/timeout.conf
# /etc/systemd/system/systemd-networkd-wait-online.service.d/timeout.conf
# Reduce the network-online wait from the default 2 minutes to 10 seconds so
# the Nightmare Loader live ISO does not hang at boot when no network
# interface becomes ready quickly (e.g. headless QEMU, offline usage).

[Service]
TimeoutStartSec=10s
